# systemd/.config/systemd/user/hyprpaper-random.service
[Unit]
Description=Cambia wallpaper aleatorio con hyprpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprpaper-random
